# patches/RoboCon2026_p5/etc/systemd/system/robot_usb_logs.service
[Unit]
Description=Log handler to collect logs onto a USB stick from Hopper

[Service]
ExecStart=/usr/local/bin/robot_usb_logs.sh
Type=simple
Restart=always
RestartSec=1
